# cmd/godoc/godoc.socket
[Socket]
BindToDevice=lo
ListenStream=6060

[Install]
WantedBy=sockets.target
